# sniproxy.ini
[program:sniproxy]
command=/usr/bin/sniproxy -f -c /etc/sniproxy.conf
