# sniproxy.ini
[program:sniproxy]
command=/usr/sbin/sniproxy -c /etc/sniproxy.conf
